# launchd/com.groove.charts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily chart scraper – runs at 07:00 every day.

  Load:    launchctl load -w ~/Library/LaunchAgents/com.groove.charts.plist
  Unload:  launchctl unload ~/Library/LaunchAgents/com.groove.charts.plist
  Run now: launchctl start com.groove.charts
  Logs:    tail -f /Volumes/Music/groove/logs/charts.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.groove.charts</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>
      export PATH="$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:/opt/homebrew/bin:$PATH";
      cd "__GROOVE_REPO__";
      exec uv run groove scrape charts
    </string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Volumes/Music/groove/logs/charts.log</string>

  <key>StandardErrorPath</key>
  <string>/Volumes/Music/groove/logs/charts-error.log</string>
</dict>
</plist>
